# deploy/aperture-booking-scheduler.service
[Unit]
Description=Aperture Booking Background Scheduler
After=network.target aperture-booking.service
Requires=aperture-booking.service

[Service]
Type=simple
User=aperture-booking
Group=aperture-booking
WorkingDirectory=/opt/aperture-booking
ExecStart=/opt/aperture-booking/venv/bin/python manage.py scheduler
Restart=on-failure
RestartSec=10

# Environment variables
Environment=DJANGO_SETTINGS_MODULE=aperture_booking.settings_production
Environment=PYTHONPATH=/opt/aperture-booking
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aperture-booking/media /opt/aperture-booking/logs /opt/aperture-booking/backups

[Install]
WantedBy=multi-user.target