# deploy/supervisord.conf
# Supervisor configuration for Aperture Booking Docker deployment

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Nginx web server
[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx/nginx.log
stderr_logfile=/var/log/nginx/nginx.error.log

# Gunicorn Django application
[program:gunicorn]
command=/opt/venv/bin/gunicorn aperture_booking.wsgi:application --config /app/gunicorn.conf.py
directory=/app
user=app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/gunicorn.log
stderr_logfile=/app/logs/gunicorn.error.log
environment=DJANGO_SETTINGS_MODULE="aperture_booking.settings_production"

# Background scheduler
[program:scheduler]
command=/opt/venv/bin/python manage.py scheduler
directory=/app
user=app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/scheduler.log
stderr_logfile=/app/logs/scheduler.error.log
environment=DJANGO_SETTINGS_MODULE="aperture_booking.settings_production"